# connman/src/connman.service.in
[Unit]
Description=Connection service
Requires=dbus.socket
After=dbus.socket android-system.service
Before=network.target remote-fs-pre.target
Wants=network.target remote-fs-pre.target

[Service]
Type=notify
Restart=always
EnvironmentFile=/etc/connman/environment.conf
ExecStart=@sbindir@/connmand -n $EXTRA_PARAM
StandardOutput=null

[Install]
WantedBy=multi-user.target
